# target/linux/ib42x0/base-files/etc/uci-defaults/migrate-leds
#!/bin/sh
sed -i -e "s/nas4220b:orange:hdd/nas42x0b:orange:hdd/g"  /etc/config/system
sed -i -e "s/nas4220b:green:os/nas42x0b:green:os/g"      /etc/config/system
