# target/linux/ib42x0/base-files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
ttyS0::askfirst:/sbin/getty 19200 ttyS0 xterm
